# templates/beat.ini
[program:{{ beat }}]
user={{ beat_user }}
command={{ beat_install_dir }}/{{ beat_version }}/{{beat}} -c {{ beat_install_dir }}/{{ beat_version }}/{{ beat }}.yml --path.data {{ beat_data_dir }} --path.logs {{ beat_log_dir }}
autostart=true
autorestart=true
stdout_logfile={{ beat_log_dir }}/{{ beat }}.out.log
stderr_logfile={{ beat_log_dir }}/{{ beat }}.err.log
